# debian/nsq-admin.upstart
description "Realtime distributed messaging (nsqadmin)"
author "Misja Hoebe <misja.hoebe@gmail.com>"

start on stopped rc RUNLEVEL=[2345]
stop on starting rc RUNLEVEL=[016]

console log
respawn

setuid nsq
setgid nsq

exec /usr/bin/nsqadmin -config /etc/nsq/nsqadmin.conf
